# templates/systemd/code-intel-mcp@.service
[Unit]
Description=code-intel MCP server (instance %i)
After=network.target

[Service]
Type=simple
# Cap glibc per-thread malloc arenas BEFORE Python starts. The code-intel CLI
# entry point has a self-reexec bootstrap that also sets this, but pre-setting
# here is belt-and-suspenders (saves one fork+exec on every restart and keeps
# the unit aligned with code-intel-watcher@.service which does NOT have a
# cli.py bootstrap).
Environment="MALLOC_ARENA_MAX=2"
# {{ code_intel_bin }} is the absolute path of the code-intel CLI at install
# time. Avoids systemd user-session PATH issues (no shims, no `code-intel`).
# The instance name (%i) maps to a manifest at
# ~/.config/code-intel/projects/%i.toml which carries the target path.
ExecStart=/bin/sh -c 'TARGET=$(awk -F\\\" "/^target/ {print \\$2}" %h/.config/code-intel/projects/%i.toml); exec {{ code_intel_bin }} serve --target "$TARGET" --stdio'
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
