# templates/systemd/code-intel-watcher@.service
[Unit]
Description=code-intel file watcher (instance %i)
After=network.target

[Service]
Type=simple
# {{ python_bin }} is the absolute path of the code-intel interpreter at install
# time. Avoids systemd user-session PATH issues (no pyenv shims, no `python`).
# The instance name (%i) maps to a manifest at
# ~/.config/code-intel/projects/%i.toml carrying the target path.
ExecStart=/bin/sh -c 'TARGET=$(awk -F\\\" "/^target/ {print \\$2}" %h/.config/code-intel/projects/%i.toml); exec {{ python_bin }} -m code_intel.watcher "$TARGET"'
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
