# templates/systemd/code-intel-watcher@.service
[Unit]
Description=code-intel file watcher (instance %i)
After=network.target

[Service]
Type=simple
# Cap glibc per-thread malloc arenas BEFORE Python starts. Belt-and-suspenders:
# v0.1.7's cli.py self-reexec covers `code-intel ...` invocations, but this
# unit launches `python -m code_intel.watcher` directly (no cli.py bootstrap),
# so the env var must be set here to plateau RSS ~220MB instead of ~589MB on
# multi-thread embedder workloads.
Environment="MALLOC_ARENA_MAX=2"
# {{ python_bin }} is the absolute path of the code-intel interpreter at install
# time. Avoids systemd user-session PATH issues (no pyenv shims, no `python`).
# The instance name (%i) maps to a manifest at
# ~/.config/code-intel/projects/%i.toml carrying the target path.
ExecStart=/bin/sh -c 'TARGET=$(awk -F\\\" "/^target/ {print \\$2}" %h/.config/code-intel/projects/%i.toml); exec {{ python_bin }} -m code_intel.watcher "$TARGET"'
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
